# .docker/php/supervisor.d/appname-worker.conf
[supervisord]
user=root

[program:appname_horizon]
process_name=%(program_name)s
command=/usr/bin/php /var/www/html/artisan horizon
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=600
